# base/php72/super.ini
[supervisord] 
nodaemon=true

[program:nginx] 
command=nginx
startretries =1

[program:php-fpm] 
command=php-fpm
startretries = 1
