# Joomla/3.8.0/php-fpm/super.ini
[supervisord] 
nodaemon=true

[program:mongo] 
command=mongod
startretries =1

[program:php-fpm] 
command=php-fpm
startretries = 1
