# base/php56/super.ini
[supervisord] 
nodaemon=true


[program:nginx]
command=nginx
numprocs=1
autostart=true
autorestart=false
startsecs=0


[program:php-fpm]
command=php-fpm5 
numprocs=1
autostart=true
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/php_fpmstdout
stdout_logfile_maxbytes=0