# base/php5.6.10/super.ini
[supervisord] 
nodaemon=true


[program:nginx]
command=nginx
numprocs=1
autostart=true
autorestart=false
startsecs=0


[program:php-fpm]
command=/usr/local/php5/sbin/php-fpm
numprocs=1
autostart=true
autorestart=false
startsecs=0

